# jcx/openrc
#!/sbin/openrc-run

name="firstboot"
command=/jcx-firstboot

depend() {
	need hostname
}
