# jcx/openrc
#!/sbin/openrc-run

name="JCX Boot Service"
command=/jcx/services/boot/runner

depend() {
	need hostname
}
