# jcx/openrc
#!/sbin/openrc-run

name="JCX First Boot Service"
command=/jcx/services/firstboot/runner

depend() {
	need hostname
}
