# jcx/sysroot/services/boot/openrc
#!/sbin/openrc-run

name="JCX Boot Service"

depend() {
	need hostname
}

start() {
	. /jcx/services/boot/runner
}